# src/platform/init/acpid.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# acpid
#
# when login-prompt-ready is received, starts acpi daemon.

start on login-prompt-ready

exec /etc/init.d/acpid start
