# src/platform/init/acpid.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# acpid
#
# when dump-boot-stats has run, starts acpi daemon.
# when halting, stop acpi daemon to keep it from responding to another event.

start on stopped dump-boot-stats
stop on starting halt

respawn
expect fork
exec /usr/sbin/acpid --confdir /etc/acpi/events --socketfile /var/run/acpid.socket
