# src/platform/init/bootchart.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Collects boot data and optionally creates a chart"
author        "chromium-os-dev@googlegroups.com"

start on stopped startup
stop on login-prompt-ready

pre-start script

BC_RUN="/var/run/bootchart"
mkdir -p "$BC_RUN"

end script

exec /lib/bootchart/collector 25 /var/run/bootchart

pre-stop exec sleep 10

post-stop script

BC_RUN="/var/run/bootchart"
BC_LOG="/var/log/bootchart"
BC_DATA="${BC_LOG}/boot-$(date +%Y%m%d-%H%M%S).tgz"

mkdir -p "$BC_LOG"
/lib/bootchart/gather "$BC_DATA" "$BC_RUN"
/usr/bin/bootchart -f png -o "$BC_LOG" "$BC_DATA"

rm -rf "$BC_RUN"

end script
