# src/platform/init/connman.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# connman
#
# starts connman once dbus has been started.

start on started wpasupplicant

script
  echo "`date` - starting connman" >> /var/log/connman.log
  /usr/sbin/connmand -n
end script
