# src/platform/init/console-kit.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Tracks user sessions and provides session cookies"
author        "chromium-os-dev@googlegroups.com"

start on started dbus
stop on starting halt or starting reboot

respawn
expect fork
exec /usr/sbin/console-kit-daemon
