# src/platform/init/cryptohomed.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# cryptohomed
#
# Starts the cryptohome daemon, which handles mounting and unmounting users'
# encrypted home directories.  Also supports offline login checks.
#

start on started dbus
stop on starting reboot or starting halt
respawn

expect fork
exec /usr/sbin/cryptohomed --noclose
