# src/platform/init/dbus.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# dbus
#
# Does some required setup and then starts the D-Bus system.
#

start on runlevel 2
respawn

# Sets up the D-Bus environment.
pre-start script
  mkdir -p /var/run/dbus
  chown messagebus.messagebus /var/run/dbus
  # Create a machine id in /var/lib/dbus
  /bin/dbus-uuidgen --ensure
end script

script
  echo "`date` - starting dbus" >> /var/log/dbus.log
  su messagebus -s /bin/sh -c 'exec /bin/dbus-daemon --system --nofork'
end script

# The post-stop script will clean up after d-dbus. It's better to put this
# here rather than in pre-start since post-stop is always guaranteed to run
# while pre-start can exit pre-maturely.
post-stop script
  rm -rf /var/run/dbus
end script
