# src/platform/init/dbus-session.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# dbus-session
#
# Sets up a login-specific DBus session.  This should be started and stopped
# for each login/logout to ensure that each logged in user has their own
# DBus session daemon.

# We only start after the ui indicates a user has logged in. On logout, ui
# is stopped.
start on start-user-session
stop on stopping ui
respawn

env USERNAME=chronos
env MINIJAIL_FLAGS="--use-capabilities"
# Must be exported for dbus-daemon-session
env DBUS_SESSION_DIR=/var/run/dbus-session

# Create a user-specific dbus session directory.
pre-start script
  mkdir -p "${DBUS_SESSION_DIR}"
  chown -R ${USERNAME}.root "${DBUS_SESSION_DIR}"
end script

script
USER_ID="$(/usr/bin/id -u ${USERNAME})"
exec /sbin/minijail ${MINIJAIL_FLAGS} --uid=${USER_ID} /sbin/dbus-daemon-session
end script

# Clean up our run directory
post-stop script
  rm -rf "${DBUS_SESSION_DIR}"
end script
