# src/platform/init/dbus-wait.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# dbus-wait
#
# this task just waits until dbus is up and emits a dbus-started task.

start on started dbus
task

script
  # Make sure D-Bus is initialized
  for i in 1 2 3 4 5 6 7 8 9 10
  do
    if [ -e /var/run/dbus/system_bus_socket ]; then
      echo "`date` - dbus is up" >> /var/log/dbus.log
      break
    fi
    echo "`date` - waiting on dbus" >> /var/log/dbus.log
    sleep .3
  done

  initctl emit dbus-started
end script
