# src/platform/init/dump-boot-stats.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# dump-boot-stats
#
# when login-prompt-ready is received, grabs the current uptime and
# disk stats

start on login-prompt-ready

# This is run-once rather than a service.
task

script
  cat /proc/uptime > /tmp/uptime-login-prompt-ready
  cat /sys/block/sda/stat > /tmp/disk-login-prompt-ready
end script
