# src/platform/init/dump-boot-stats.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# dump-boot-stats
#
# when login-prompt-ready is received, grabs the current uptime (in seconds)
# and stores it to a file in tmp (in ms)

start on login-prompt-ready

# This is run-once rather than a service.
task

exec cat /proc/uptime | awk '{print $1 * 1000}' > /tmp/boot-time
