# src/platform/init/googleguest.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

#
#
# Automatically connects to GoogleGuest or any other open network
# once connman comes up.
#
# This process is skipped if a user has succesfully
# logged into the machine since Connman will attempt to reconnect to any
# networks that it has connected to previously.

start on started connman
console output

# Tells upstart not to manage this script as a service so that it
# is only run once.
task

script
  date > /var/log/gguest.log
  /usr/local/connman/test/start-scanning
  if [ ! -f /var/cache/google_offline_login_cache.txt ]
  then
    /opt/google/controlpanel/connect-to-open-wifi >> /var/log/gguest.log
  else
    echo "Skipping network configuration" >> /var/log/gguest.log
  fi
end script
