# src/platform/init/halt.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Halt the system
#
# TODO: For now we just execute what is in rc0 for the most part.

description   "Halt (shutdown) the system"
author        "chromium-os-dev@googlegroups.com"

start on runlevel 0
stop on runlevel [!0]

task
console output

script

# TERM chrome before the user's home dir is unmounted by the murder of slim.
touch /tmp/disable_chrome_restart
pkill -TERM "^chrome$"
/sbin/hwclock --systohc --utc
/etc/init.d/slim stop
/etc/init.d/alsa-utils stop
/etc/init.d/sendsigs stop
/etc/init.d/umountfs stop
/etc/init.d/halt stop

end script

