# src/platform/init/ibus.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run ibus for the given logged in user"
author        "chromium-os-dev@googlegroups.com"

# The dbus-session should be started and stopped for each user session.
start on started dbus-session
stop on stopping dbus-session
respawn

env USERNAME=chronos
env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0
env MINIJAIL_FLAGS="--use-capabilities"
# Setup by dbus-session.conf
env DBUS_SESSION_BUS_ADDRESS_FILE=/var/run/dbus-session/address

script
# Get the numeric uid
USER_ID="$(/usr/bin/id -u ${USERNAME})"
# ibus checks this env variable to get the UID of the calling user.  It is
# meant for cases where userhelper is used, but it seems to be needed for us.
export USERHELPER_UID=${USER_ID}
# We import the DBus bus address for our given session (for the user login).
# GConf2, if used, may supercede this need.
export DBUS_SESSION_BUS_ADDRESS="$(/bin/cat $DBUS_SESSION_BUS_ADDRESS_FILE)"

exec /sbin/minijail ${MINIJAIL_FLAGS} --uid=${USER_ID} -- /usr/bin/ibus-daemon
end script
