# src/platform/init/ibus.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run ibus for the given logged in user"
author        "chromium-os-dev@googlegroups.com"

# start-user-session is only emitted from the new session_manager code.
# chromeos-xsession (the old way of starting the WM) is not run on that
# code pathway.
start on start-user-session
stop on starting halt or starting reboot # at least, for now
# stop on stop-user-session

env USER_ID=1000
env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0

respawn
expect fork

script
# this is some weird env var that ibus-daemon checks to determine what
# user it's running as
export USERHELPER_UID=${USER_ID}
exec /sbin/minijail --uid=${USER_ID} -- /usr/bin/ibus-daemon
end script
