# src/platform/init/klog.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
#

description   "Pipes data from kmsg to the syslog daemon"
author        "chromium-os-dev@googlegroups.com"

start on started syslog
stop on stopping syslog

# We use 'cat' because we want to read line buffered i/o from /proc/kmsg and
# pass it on to the syslog daemon through the kmsg pipe.
respawn
exec /bin/cat < /proc/kmsg > /var/run/rsyslog/kmsg
