# src/platform/init/metrics_daemon.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# metrics_daemon

# We'd like to catch the early network stats.
start on started connman

respawn
expect fork
exec /usr/bin/metrics_daemon
