# src/platform/init/netdebug.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# netdebug
#
# enables network debugging once connman has been started.

start on started connman

# Tells upstart not to manage this script as a service so that it
# is only run once.
task

script
  # enable connman wifi debug msgs
  /usr/local/connman/test/debug wifi
  # set wpa_supplicant debug level to 1, disable timestamps, disable show_keys
  dbus-send --system --dest=fi.epitest.hostap.WPASupplicant --print-reply \
    /fi/epitest/hostap/WPASupplicant \
    fi.epitest.hostap.WPASupplicant.setDebugParams \
    int32:1 boolean:false boolean:false
  # add other stuff here
end script
