# src/platform/init/ntp.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# ntp
#
# when udev-settled is received, starts ntp daemon.

start on udev-settled

script
  echo "`date` - starting ntp" >> /var/log/ntp.log

  # Kill any previous instances of ntpd. Running pkill in a new shell because
  # if pkill exits with an error state, this whole script will stop running.
  sh -c "pkill -9 ntpd; exit 0;"

  /etc/init.d/ntp start >> /var/log/ntp.log
  /etc/init.d/ntp status >> /var/log/ntp.log
end script
