# src/platform/init/ntp.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# ntp
#
# Starts the ntp daemon to set the system time.

start on started udev

respawn
expect fork
exec /usr/sbin/ntpd -g -u ntp:ntp

