# src/platform/init/omaha-tracker.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the omaha ID tracker for the given logged in user"
author        "chromium-os-dev@chromium.org"

# as the script /usr/sbin/omaha_tracker.sh only has executable content on
# Google Chrome OS builds, this job is a no-op on chromium OS builds.

# start-user-session is only emitted from the new session_manager code.
# chromeos-xsession (the old way of starting the WM) is not run on that
# code pathway.
start on start-user-session
stop on stopping ui

# sadly, these can't reference each other.
env USER=chronos
env DATA_DIR=/home/chronos
env HOME=/home/chronos/user

respawn
normal exit 0 TERM

script
OMAHA_TRACKER=/usr/sbin/omaha_tracker.sh
if [ $(stat -c%s ${OMAHA_TRACKER}) = "0" ]; then
  exit 0
else
  exec /sbin/minijail --uid=$(/usr/bin/id -u ${USER}) -- ${OMAHA_TRACKER}
fi
end script
