# src/platform/init/openssh-server.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

start on started connman
respawn
expect fork
exec /usr/sbin/sshd



pre-start script 

mkdir -p /mnt/stateful_partition/etc/ssh

if [ ! -f /mnt/stateful_partition/etc/ssh/ssh_host_rsa_key ] ; then
   echo Generating v2 RSA key
   ssh-keygen -q -f /mnt/stateful_partition/etc/ssh/ssh_host_rsa_key -N '' -t rsa
fi


if [ ! -f /mnt/stateful_partition/etc/ssh/ssh_host_dsa_key ] ; then
   echo Generating v2 DSA key
   ssh-keygen -q -f /mnt/stateful_partition/etc/ssh/ssh_host_dsa_key -N '' -t dsa
fi

end script 


