# src/platform/init/preload-chrome.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# preload-chrome
#
# when login-prompt-ready is received, dumps all the files touched by chromium
# during startup out to /dev/null to prewarm the disk cache.

start on login-prompt-ready

# This is run-once rather than a service.
task

script
{
set +e
# This list generated from an strace of Chromium startup.
FILES="/opt/google/chrome/chrome /etc/fonts/conf.d/10-antialias.conf /etc/fonts/conf.d/10-hinting-slight.conf /etc/fonts/conf.d/10-hinting.conf /etc/fonts/conf.d/11-lcd-filter-lcddefault.conf /etc/fonts/conf.d/20-fix-globaladvance.conf /etc/fonts/conf.d/20-unhint-small-vera.conf /etc/fonts/conf.d/30-defoma.conf /etc/fonts/conf.d/30-metric-aliases.conf /etc/fonts/conf.d/30-urw-aliases.conf /etc/fonts/conf.d/40-nonlatin.conf /etc/fonts/conf.d/45-latin.conf /etc/fonts/conf.d/49-sansserif.conf /etc/fonts/conf.d/50-user.conf /etc/fonts/conf.d/51-local.conf /etc/fonts/conf.d/53-monospace-lcd-filter.conf /etc/fonts/conf.d/60-latin.conf /etc/fonts/conf.d/65-fonts-persian.conf /etc/fonts/conf.d/65-nonlatin.conf /etc/fonts/conf.d/69-unifont.conf /etc/fonts/conf.d/70-no-bitmaps.conf /etc/fonts/conf.d/80-delicious.conf /etc/fonts/conf.d/90-synthetic.conf /etc/fonts/fonts.conf /etc/gai.conf /etc/gtk-2.0/gdk-pixbuf.loaders /etc/gtk-2.0/gtk.immodules /etc/gtk-2.0/gtkrc /etc/host.conf /etc/hosts /etc/ld.so.cache /etc/localtime /etc/lsb-release /etc/nsswitch.conf /etc/pango/pango.modules /etc/pango/pangorc /etc/passwd /etc/resolv.conf /etc/selinux/config /lib/cmov/libnss3.so.1d /lib/i686/cmov/libnspr4.so.0d /lib/i686/cmov/libnss3.so.1d /lib/i686/cmov/libnssutil3.so.1d /lib/i686/cmov/libplc4.so.0d /lib/i686/cmov/libplds4.so.0d /lib/i686/cmov/libsmime3.so.1d /lib/i686/cmov/libssl3.so.1d /lib/i686/libnspr4.so.0d /lib/i686/libnss3.so.1d /lib/i686/libnssutil3.so.1d /lib/i686/libplc4.so.0d /lib/i686/libplds4.so.0d /lib/i686/libsmime3.so.1d /lib/i686/libssl3.so.1d /lib/i686/sse2/cmov/libnss3.so.1d /lib/i686/sse2/libnss3.so.1d /lib/libdbus-1.so.3 /lib/libgcc_s.so.1 /lib/libnspr4.so.0d /lib/libnss3.so.1d /lib/libnssutil3.so.1d /lib/libpcre.so.3 /lib/libplc4.so.0d /lib/libplds4.so.0d /lib/libselinux.so.1 /lib/libsmime3.so.1d /lib/libssl3.so.1d /lib/libz.so.1 /lib/sse2/cmov/libnss3.so.1d /lib/sse2/libnss3.so.1d /lib/tls/cmov/libnss3.so.1d /lib/tls/i686/cmov/libc.so.6 /lib/tls/i686/cmov/libdl.so.2 /lib/tls/i686/cmov/libm.so.6 /lib/tls/i686/cmov/libnsl.so.1 /lib/tls/i686/cmov/libnspr4.so.0d /lib/tls/i686/cmov/libnss3.so.1d /lib/tls/i686/cmov/libnss_compat.so.2 /lib/tls/i686/cmov/libnss_dns.so.2 /lib/tls/i686/cmov/libnss_files.so.2 /lib/tls/i686/cmov/libnss_nis.so.2 /lib/tls/i686/cmov/libnssutil3.so.1d /lib/tls/i686/cmov/libplc4.so.0d /lib/tls/i686/cmov/libplds4.so.0d /lib/tls/i686/cmov/libpthread.so.0 /lib/tls/i686/cmov/libresolv.so.2 /lib/tls/i686/cmov/librt.so.1 /lib/tls/i686/cmov/libsmime3.so.1d /lib/tls/i686/cmov/libssl3.so.1d /lib/tls/i686/cmov/libutil.so.1 /lib/tls/i686/libnspr4.so.0d /lib/tls/i686/libnss3.so.1d /lib/tls/i686/libnssutil3.so.1d /lib/tls/i686/libplc4.so.0d /lib/tls/i686/libplds4.so.0d /lib/tls/i686/libsmime3.so.1d /lib/tls/i686/libssl3.so.1d /lib/tls/i686/sse2/cmov/libnss3.so.1d /lib/tls/i686/sse2/libnss3.so.1d /lib/tls/libnspr4.so.0d /lib/tls/libnss3.so.1d /lib/tls/libnssutil3.so.1d /lib/tls/libplc4.so.0d /lib/tls/libplds4.so.0d /lib/tls/libsmime3.so.1d /lib/tls/libssl3.so.1d /lib/tls/sse2/cmov/libnss3.so.1d /lib/tls/sse2/libnss3.so.1d /opt/google/chrome/chrome.pak /opt/google/chrome/chromeos/libcros.so /opt/google/chrome/libavcodec.so.52 /opt/google/chrome/libavformat.so.52 /opt/google/chrome/libavutil.so.50 /opt/google/chrome/locales/en-US.pak /opt/google/chrome/themes/default.pak /usr/X11R6/lib/X11/icons/default/cursors/08e8e1c95fe2fc01f976f1e063a24ccd /usr/X11R6/lib/X11/icons/default/cursors/sb_v_double_arrow /usr/X11R6/lib/X11/icons/default/cursors/xterm /usr/bin/chromeos-chrome /usr/bin/commands.so /usr/bin/commandsmodule.so /usr/bin/functools.so /usr/bin/functoolsmodule.so /usr/bin/gettext.so /usr/bin/gettextmodule.so /usr/bin/locale.so /usr/bin/localemodule.so /usr/bin/lsb_release /usr/bin/lsb_release.so /usr/bin/lsb_releasemodule.so /usr/bin/optparse.so /usr/bin/optparsemodule.so /usr/bin/org.so /usr/bin/orgmodule.so /usr/bin/re.so /usr/bin/remodule.so /usr/bin/sre_compile.so /usr/bin/sre_compilemodule.so /usr/bin/sre_constants.so /usr/bin/sre_constantsmodule.so /usr/bin/sre_parse.so /usr/bin/sre_parsemodule.so /usr/bin/string.so /usr/bin/stringmodule.so /usr/bin/struct.so /usr/bin/structmodule.so /usr/bin/textwrap.so /usr/bin/textwrapmodule.so /usr/lib/charset.alias /usr/lib/cmov/libnss3.so.1d /usr/lib/gconv/gconv-modules.cache /usr/lib/gtk-2.0/2.10.0/immodule-files.d/libgtk2.0-0.immodules /usr/lib/gtk-2.0/2.10.0/loader-files.d/libgtk2.0-0.loaders /usr/lib/gtk-2.0/2.10.0/loaders/libpixbufloader-png.so /usr/lib/i686/cmov/libnss3.so.1d /usr/lib/i686/libnss3.so.1d /usr/lib/i686/sse2/cmov/libnss3.so.1d /usr/lib/i686/sse2/libnss3.so.1d /usr/lib/libORBit-2.so.0 /usr/lib/libX11.so.6 /usr/lib/libXau.so.6 /usr/lib/libXcomposite.so.1 /usr/lib/libXcursor.so.1 /usr/lib/libXdamage.so.1 /usr/lib/libXdmcp.so.6 /usr/lib/libXext.so.6 /usr/lib/libXfixes.so.3 /usr/lib/libXi.so.6 /usr/lib/libXinerama.so.1 /usr/lib/libXrandr.so.2 /usr/lib/libXrender.so.1 /usr/lib/libasound.so.2 /usr/lib/libatk-1.0.so.0 /usr/lib/libcairo.so.2 /usr/lib/libdbus-glib-1.so.2 /usr/lib/libdirect-1.2.so.0 /usr/lib/libdirectfb-1.2.so.0 /usr/lib/libexpat.so.1 /usr/lib/libfontconfig.so.1 /usr/lib/libfreetype.so.6 /usr/lib/libfusion-1.2.so.0 /usr/lib/libgconf-2.so.4 /usr/lib/libgdk-x11-2.0.so.0 /usr/lib/libgdk_pixbuf-2.0.so.0 /usr/lib/libgflags.so.0 /usr/lib/libgio-2.0.so.0 /usr/lib/libglib-2.0.so.0 /usr/lib/libglog.so.0 /usr/lib/libgmodule-2.0.so.0 /usr/lib/libgobject-2.0.so.0 /usr/lib/libgthread-2.0.so.0 /usr/lib/libgtk-x11-2.0.so.0 /usr/lib/libnspr4.so.0d /usr/lib/libnss3.so.1d /usr/lib/libnssutil3.so.1d /usr/lib/libpango-1.0.so.0 /usr/lib/libpangocairo-1.0.so.0 /usr/lib/libpangoft2-1.0.so.0 /usr/lib/libpixman-1.so.0 /usr/lib/libplc4.so.0d /usr/lib/libplds4.so.0d /usr/lib/libpng12.so.0 /usr/lib/libsmime3.so.1d /usr/lib/libssl3.so.1d /usr/lib/libstdc++.so.6 /usr/lib/libxcb-render-util.so.0 /usr/lib/libxcb-render.so.0 /usr/lib/libxcb.so.1 /usr/lib/pango/1.6.0/module-files.d/libpango1.0-0.modules /usr/lib/pango/1.6.0/modules/pango-basic-fc.so /usr/lib/sse2/cmov/libnss3.so.1d /usr/lib/sse2/libnss3.so.1d /usr/lib/tls/cmov/libnss3.so.1d /usr/lib/tls/i686/cmov/libnss3.so.1d /usr/lib/tls/i686/libnss3.so.1d /usr/lib/tls/i686/sse2/cmov/libnss3.so.1d /usr/lib/tls/i686/sse2/libnss3.so.1d /usr/lib/tls/libnss3.so.1d /usr/lib/tls/sse2/cmov/libnss3.so.1d /usr/lib/tls/sse2/libnss3.so.1d /usr/share/X11/XKeysymDB /usr/share/X11/locale/C/XLC_LOCALE /usr/share/X11/locale/locale.alias /usr/share/X11/locale/locale.dir /usr/share/fonts/truetype/ttf-droid/DroidSans-Bold.ttf /usr/share/fonts/truetype/ttf-droid/DroidSans.ttf /usr/share/fonts/truetype/ttf-droid/DroidSansFallback.ttf /usr/share/icons/chromeos-cursors/cursors/08e8e1c95fe2fc01f976f1e063a24ccd /usr/share/icons/chromeos-cursors/cursors/sb_v_double_arrow /usr/share/icons/chromeos-cursors/cursors/xterm /usr/share/icons/default/cursors/08e8e1c95fe2fc01f976f1e063a24ccd /usr/share/icons/default/cursors/sb_v_double_arrow /usr/share/icons/default/cursors/xterm /usr/share/icons/default/index.theme /usr/share/locale/locale.alias /usr/share/pixmaps/default/cursors/08e8e1c95fe2fc01f976f1e063a24ccd /usr/share/pixmaps/default/cursors/sb_v_double_arrow /usr/share/pixmaps/default/cursors/xterm /usr/share/themes/Raleigh/gtk-2.0/gtkrc external_extensions.json"

for file in ${FILES}
do
  dd if=${file} of=/dev/null
done
set -e
} &
end script
