# src/platform/init/procps.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# procps
#
# when login-prompt-ready is received, runs procps which sets sysctls.

start on login-prompt-ready

# This is run-once rather than a service.
task

exec /etc/init.d/procps start
