# src/platform/init/reboot.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Reboot the system"
author        "chromium-os-dev@googlegroups.com"

start on runlevel 6
stop on runlevel [!6]

console output

task
exec /etc/init.d/chromeos_shutdown
post-stop exec /sbin/halt --force
