# src/platform/init/reboot.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Reboot the system
#
# TODO: For now we just execute what is in rc6 for the most part.

description   "Reboot the system"
author        "chromium-os-dev@googlegroups.com"

start on runlevel 6
stop on runlevel [!6]

task
console output

script

/sbin/hwclock --systohc --utc
/etc/init.d/slim stop
/etc/init.d/alsa-utils stop
/etc/init.d/sendsigs stop
/etc/init.d/umountfs stop
/etc/init.d/reboot stop

end script

