# src/platform/init/screen-locker.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the screen locker for the given logged in user"
author        "chromium-os-dev@chromium.org"

# start-user-session is only emitted from the new session_manager code.
# chromeos-xsession (the old way of starting the WM) is not run on that
# code pathway.
start on start-user-session
stop on stopping ui

# sadly, these can't reference each other.
env USER=chronos
env DATA_DIR=/home/chronos
env HOME=/home/chronos/user
env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0

respawn

exec /sbin/minijail --uid=$(/usr/bin/id -u ${USER}) -- \
  /usr/bin/xscreensaver -no-splash
