# src/platform/init/software-update.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# software-update
#
# when login-prompt-ready is received, start up software update check.

start on login-prompt-ready

script
  DEVKIT_URL=$(grep ^CHROMEOS_DEVSERVER /etc/lsb-release | cut -d = -f 2-) 
  if [ "x" = "x$DEVKIT_URL" ]
  then
    /etc/init.d/software_update.sh start
  fi
end script
