# src/platform/init/software-update.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# description "System software update service"
# author "chromium-os-dev@googlegroups.com"
#
# when dump-boot-stats has run, start up software update check.

start on stopped dump-boot-stats

respawn

exec /opt/google/memento_updater/software_update.sh 1800  # seconds = 30 min
