# src/platform/init/stop-bootchart.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# stop-bootchart
#
# when login-prompt-ready is received, waits a bit and then stops bootchart.

start on login-prompt-ready

# Indicate that this is run-once rather than a service.
task

script
if [ -x /lib/bootchart/collector ]
then
  {
    sleep 7
    /etc/init.d/stop-bootchart start
  } &
fi
end script
