# src/platform/init/syslog.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
#

description   "Syslog daemon"
author        "chromium-os-dev@googlegroups.com"

start on runlevel 2
stop on starting halt or starting reboot

pre-start script

SYSLOG_RUN="/var/run/rsyslog"
KMSG_PIPE="${SYSLOG_RUN}/kmsg"

mkdir -p "$SYSLOG_RUN"
test -p "$KMSG_PIPE" || mkfifo -m 600 "$KMSG_PIPE"
chown -R syslog:syslog "$SYSLOG_RUN"

end script

respawn
expect fork
exec /usr/sbin/rsyslogd -c4
