# src/platform/init/syslog.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
#
start on runlevel 2

respawn

pre-start script
  mkdir -p /var/run/rsyslog
  mkfifo -m 600 /var/run/rsyslog/kmsg
  chown -R syslog:syslog /var/run/rsyslog

  /sbin/start-stop-daemon --background --start --quiet \
    --exec /bin/dd -- if=/proc/kmsg of=/var/run/rsyslog/kmsg
end script

expect fork
exec /usr/sbin/rsyslogd -c4
