# src/platform/init/syslog.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# syslog
#
# when login-prompt-ready is received, starts syslog and kernel log daemons.

start on login-prompt-ready

script
{
  /etc/init.d/rsyslog start
  /etc/init.d/bootlogs.sh start
} &
end script
