# src/platform/init/udev.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev
#
# when login-prompt-ready is received, starts udev.

start on login-prompt-ready

script
{
  start-stop-daemon --start --quiet --exec /sbin/udevd -- --daemon
  /sbin/udevadm trigger
  /sbin/udevadm settle

  initctl emit udev-settled
} &
end script
