# src/platform/init/udev.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev
#
# when login-prompt-ready is received, starts udev.

start on login-prompt-ready

respawn
expect fork
exec /sbin/udevd --daemon

post-start script
  /sbin/udevadm trigger
  /sbin/udevadm settle
end script
