# src/platform/init/ui.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Chrome OS user interface"
author        "chromium-os-dev@googlegroups.com"

start on stopping startup
stop on starting halt or starting reboot

respawn

pre-start script

X_SOCKET_DIR=/tmp/.X11-unix
X_ICE_DIR=/tmp/.ICE-unix
mkdir -p "$X_SOCKET_DIR" "$X_ICE_DIR"
chown root:root "$X_SOCKET_DIR" "$X_ICE_DIR"
chmod 1777 "$X_SOCKET_DIR" "$X_ICE_DIR"

end script

# To switch to experimental boot-to-chrome functionality, uncomment the
# following and comment out the lines below it.
#exec /etc/init.d/session_manager.sh
expect fork
exec /usr/bin/slim -d

pre-stop script

# Shutdown is best-effort, so don't let errors kill this part way through.
set +e

# Tries to cleanly terminate a process and waits for it to exit.
# If that fails it will kill the process.
term_process() {
  local process="$1"

  pkill -o "$process"
  for i in 1 2 3 4 5 6 7 8 9 10 11 12 13 14 15; do
    pgrep "$process"
    if [ $? -ne 0 ] ; then
      return
    fi
    sleep .1
  done
  pkill -KILL "$process"
}
touch /tmp/disable_chrome_restart
term_process "^chrome$"

end script
