# src/platform/init/ui.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Chrome OS user interface"
author        "chromium-os-dev@googlegroups.com"

start on stopping startup
stop on starting halt or starting reboot

respawn

pre-start script

X_SOCKET_DIR=/tmp/.X11-unix
X_ICE_DIR=/tmp/.ICE-unix
mkdir -p "$X_SOCKET_DIR" "$X_ICE_DIR"
chown root:root "$X_SOCKET_DIR" "$X_ICE_DIR"
chmod 1777 "$X_SOCKET_DIR" "$X_ICE_DIR"

end script

#expect fork
#exec /usr/bin/slim -d
# To switch back to deprecated slim-based login functionality, comment out the
# following and uncomment the two lines above.
exec /sbin/session_manager_setup.sh

post-stop script

set +e
. /sbin/killers
kill_with_open_files_on /home/chronos/user
term_process "^Xorg$"
# Eventually, this will take a parameter specifying which user's dir to unmount.
exec /usr/sbin/umount.cryptohome

end script
