# src/platform/init/windowmanager.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the window manager for the given logged in user"
author        "chromium-os-dev@chromium.org"

# start-user-session is only emitted from the new session_manager code.
# chromeos-xsession (the old way of starting the WM) is not run on that
# code pathway.
start on start-user-session
stop on starting halt or starting reboot # at least, for now
# stop on stop-user-session

# sadly, these can't reference each other.
env USER=chronos
env DATA_DIR=/home/chronos
env HOME=/home/chronos/user
env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0

respawn
script

WM="/usr/bin/chromeos-wm"
IMAGES="/usr/share/chromeos-assets/images"
BACKGROUND_FILE="background_1024x600.png"
BACKGROUND="${IMAGES}/${BACKGROUND_FILE}"
LOG_DIR="${HOME}/log"
XTERM_COMMAND="/usr/bin/xterm -bg black -fg grey"

cd ${HOME}
exec /sbin/minijail --uid=$(/usr/bin/id -u ${USER}) -- \
  "${WM}"                                              \
    --hotkey_overlay_image_dir="${IMAGES}"             \
    --panel_anchor_image="${IMAGES}/panel_anchor.png"  \
    --panel_bar_image="${IMAGES}/panel_bar_bg.png"     \
    --shadow_image_dir="${IMAGES}"                     \
    --wm_background_image="${BACKGROUND}"              \
    --wm_xterm_command="${XTERM_COMMAND}"              \
    --log_dir="${LOG_DIR}"
end script
