# src/platform/init/wpasupplicant.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# wpa_supplicant
#
# starts wpa_supplicant once dbus has been started.

start on started dbus
stop on runlevel [06]
respawn

exec /sbin/wpa_supplicant -u -s
